# etc/cron.d/app
APP=vertigo
LOG=/srv/$APP/var/log
PATH=$PATH/usr/sbin:/sbin:/bin:/srv/$APP/bin/

0    */6  * * *  root push.sh >> $LOG/push.log 2>&1
*/10 * 	  * * *  root poll_twitter.sh >> $LOG/pull_twitter.log 2>&1
#0,30 *    * * * root upgrade.sh >> $LOG/upgrade.log 2>&1
