# etc/cron.d/app
APP=/srv/vertigo
PATH=/usr/bin:/bin:/usr/sbin:/sbin:/usr/local/bin

0    */6  * * *  root cd $APP; bin/push.sh >> var/log/cron/push.log
*/10 *    * * *  root cd $APP; bin/poll_twitter.sh >> var/log/cron/pull_twitter.log
#0,30 *    * * * root cd $APP; bin/upgrade.sh >> var/log/cron/upgrade.log
