# etc/cron.d/app
APP=/srv/vertigo                                                                                                                                                 LOG=/srv/vertigo/var/log/cron
PATH=$PATH:/usr/sbin:/sbin:/bin:/usr/local/bin

0    */6  * * *  root cd $APP; bin/push.sh >> $LOG/push.log
*    *    * * *  root cd $APP; bin/poll_twitter.sh; touch $LOG/pull_twitter.log
#0,30 *    * * * root upgrade.sh >> $LOG/upgrade.log
